# scripts/linux/homeserverhub.service
# Service systemd pour HomeServerHub.
#
# Le service tourne EN TANT QUE L'UTILISATEUR (User=/Group=), pour que ses
# données restent accessibles dans ~/.local/share/morfredus/HomeServerHub
# (et non dans /var/lib, réservé au root). Le script install-service.sh
# remplace __RUN_USER__ par l'utilisateur courant ; en installation manuelle,
# remplacez-le vous-même (ou ajoutez un drop-in User=/Group=).
#
# Installation manuelle :
#   sudo cp HomeServerHub /usr/local/bin/
#   sudo mkdir -p /etc/homeserverhub && sudo cp config.example.json /etc/homeserverhub/config.json
#   sed 's/__RUN_USER__/votreuser/g' scripts/linux/homeserverhub.service | sudo tee /etc/systemd/system/homeserverhub.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now homeserverhub

[Unit]
Description=HomeServerHub — socle de synchronisation offline-first
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
ExecStart=/usr/local/bin/HomeServerHub /etc/homeserverhub/config.json
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
